# crontab
#00 * * * * /etc/asterisk/scripts/saytime.sh
#10,20,30,40,50 * * * * /etc/asterisk/scripts/sayid.sh
*/10 * * * * /etc/asterisk/scripts/sayidandtime.sh
